# resources/linux/DEBIAN/prerm
#!/bin/sh

systemctl stop vpnht &> /dev/null || true
systemctl disable vpnht &> /dev/null || true