# resources/linux/vpnht.service
[Unit]
Description=VPN.ht Client Daemon

[Service]
ExecStart=/usr/bin/vpnht-service

[Install]
WantedBy=multi-user.target
